# crontab
0 2 * * * root cd /; ./certbot-auto renew --non-interactive > /var/log/letsencrypt-cron.log
