# crontab
0 2 * * * root /usr/bin/certbot renew --non-interactive > /var/log/letsencrypt-cron.log
