# crontab
0 2 * * * root cd /letsencrypt; ./letsencrypt-auto renew --non-interactive > /var/log/letsencrypt-cron.log
